# supervisor/supervisord-cloudsql.ini
[program:cloudsql-proxy]
command=/cloud_sql_proxy --dir=/cloudsql -instances=%(ENV_CLOUDSQL_INSTANCE)s -credential_file=%(ENV_GOOGLE_CLOUD_PUBSUB_KEY)s
autostart=true
autorestart=true
priority=1
startretries=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

